# deploy/cache-metrics.service
[Unit]
Description=Emit cargo cache metrics to node_exporter textfile collector
After=network-online.target

[Service]
Type=oneshot
User=ltodorov
Environment=CACHE_DIR=/home/ltodorov/rp2026/data/cargo-cache
Environment=OUT=/var/lib/prometheus/node-exporter/cargo_cache.prom
ExecStart=/home/ltodorov/rp2026/deploy/cache-metrics.sh
Nice=10
IOSchedulingClass=idle
